# dns_stats.service
[Unit]
Description=Pi-hole Visualizer Service

[Service]
Type=idle
ExecStart=/usr/bin/python3 /home/pi/pi-hole-visualizer/dns_stats.py -c alternate -r

[Install]
WantedBy=multi-user.target
